# scripts/supervisor/ir_wireless.conf
[program:ir_wireless]
command = /opt/ir_wireless/scripts/control_server.py
autostart = true
autorestart = true
stopsignal = KILL
